# acquisition/meterhub-acquisition.service
[Unit]
Description=MeterHub Acquisition Service - Modbus RTU Polling
Documentation=file:///opt/meterhub/docs/ARCHITECTURE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=meterhub
Group=meterhub

# Working directory
WorkingDirectory=/opt/meterhub/current

# Executable
ExecStart=/opt/meterhub/current/venv/bin/python -m meterhub_acq

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Resource limits (strict)
MemoryMax=42M
MemoryAccounting=true
CPUQuota=50%  # 50% of single core on RPi Zero W

# Security
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/meterhub /var/log/meterhub /etc/meterhub/secrets

# Environment
Environment="METREHUB_ENV=production"
EnvironmentFile=/etc/meterhub/config.env
EnvironmentFile=-/etc/meterhub/secrets.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meterhub-acq

# Timeout
TimeoutStopSec=30

# Watchdog (bcm2835_wdt)
WatchdogSec=60

[Install]
WantedBy=multi-user.target
